# backend/deployment/systemd/n8n.service
[Unit]
Description=n8n - Workflow Automation for capibara6-consensus
After=network.target flask-app.service

[Service]
Type=simple
User=elect
WorkingDirectory=/home/elect
Environment="N8N_HOST=0.0.0.0"
Environment="N8N_PORT=5678"
Environment="N8N_PROTOCOL=http"
Environment="NODE_ENV=production"
Environment="WEBHOOK_URL=http://localhost:5678/"
Environment="N8N_USER_FOLDER=/home/elect/.n8n"
Environment="N8N_LOG_LEVEL=info"
Environment="N8N_LOG_OUTPUT=console,file"
Environment="N8N_LOG_FILE_LOCATION=/home/elect/.n8n/logs/"
ExecStart=/usr/bin/n8n start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
